# supervisord.conf
[supervisord]
nodaemon=true

[program:nginxd]
command=/usr/local/nginx/sbin/nginx

[program:php-fpm]
command=/usr/local/sbin/php-fpm
